# src/etc/systemd/system/elf_monitor.service
[Unit]
Description=X64 ELF Monitor Service
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/bin/elf_monitor
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
